# etc/systemd/system/docker.service.d/override.conf
[Unit]
Requires=media.service
After=media.service
